# pobs/tools/dockerfile_augmentation/pobs_files/etc/services.d/pobs/run
#!/usr/bin/with-contenv sh

echo "system call monitor started"
strace -qqq -ttt -T -Z -k -f --trace=%file,%memory,%net -p $(pidof java)
